# gorum_1.0.16/DEBIAN/prerm
#!/bin/sh

systemctl stop gorum.service
systemctl disable gorum.service

exit 0
